# crontab
#put your cronjobs here
* * * * * /usr/local/bin/example-script